# dot_config/systemd/user/battery-user.service
[Unit]
Description=battery-user

[Service]
ExecStart=/home/sganon/.config/systemd/scripts/check-battery.sh
